# conf/upstart-job/worldcup2014.conf
#!upstart
#
# An upstart script for running a Node.js process as a service using
# Forever as the process monitor.
# Inspired from
# https://www.exratione.com/2013/02/nodejs-and-forever-as-a-service-simple-upstart-and-init-scripts-for-ubuntu/
# http://technosophos.com/2013/03/06/how-use-ubuntus-upstart-control-nodejs-forever.html
# http://serverfault.com/questions/506999/simple-upstart-script-for-nodejs-forever-on-ubuntu
# http://upstart.ubuntu.com/cookbook/
#
description "Upstart script for the worlcup2014 Node.js process"
author "Stéphane Fidanza <stephane.fidanza@gmail.com>"

start on filesystem and started networking
stop on shutdown

# This line is needed so that Upstart reports the pid of the Node.js process
# started by Forever rather than Forever's pid.
expect fork

# Let's make sure we always have the right directory
chdir /var/www/worldcup2014/src

script
    exec forever start --uid "worldcup2014" index.js
end script
 
pre-stop script
    exec forever stop "worldcup2014"
end script
